# weibo-monitor.service
[Unit]
Description=Weibo Hot Search Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/weibo-hotsearch
ExecStart=/opt/weibo-hotsearch/venv/bin/python weibo_hotsearch_monitor.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/weibo-hotsearch/logs/stdout.log
StandardError=append:/opt/weibo-hotsearch/logs/stderr.log

[Install]
WantedBy=multi-user.target
